# backupRotation.service
[Unit]
Description=Backup rotation logic.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backupRotation.sh
